# minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=exec
WorkingDirectory=/home/minecraft/papermc
EnvironmentFile=/home/minecraft/papermc/minecraft_java.env
ExecStart=java $JAVA_OPTS -jar /home/minecraft/papermc/server.jar nogui
ExecStop=/bin/sh -c 'echo stop > /run/minecraft.stdin'
User=minecraft
Restart=on-failure
RestartSec=5
TimeoutStopSec=120
Sockets=minecraft-stdin.socket
StandardInput=socket
StandardOutput=journal
StandardError=journal

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/minecraft/papermc
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
