# minecraft-stdin.socket
[Unit]
PartOf=minecraft.service

[Socket]
ListenFIFO=%t/minecraft.stdin
SocketMode=0660
SocketUser=minecraft
SocketGroup=minecraft